# com.jy.morning_alert.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jy.morning_alert</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/jy/Documents/Claude/Projects/개인비서/morning_alert.py</string>
    </array>

    <!-- 매일 오전 8시 실행 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/jy/Documents/Claude/Projects/개인비서/alert.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jy/Documents/Claude/Projects/개인비서/alert.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
